# systemd/atelier.service
[Unit]
Description=Atelier — Application Platform
# Le control-plane (registres apps+ports, tâches, index docs) vit désormais dans
# Postgres (atelier_meta). La supervision en dépend au démarrage — on ordonne donc
# après postgresql. Pas de Requires= : si Postgres tombe à chaud, Atelier continue
# de tourner (les écritures dégradent), mais un boot sans Postgres échoue volontairement.
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/atelier/bin/atelier
WorkingDirectory=/opt/atelier
EnvironmentFile=-/opt/atelier/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atelier

# Arrêt : SIGTERM au SEUL process principal d'abord (KillMode=mixed), pas à tout le cgroup.
# WHY : à l'arrêt, Atelier draine ses runners agent en vol (interrupt du tour + EOF stdin →
# le SDK flush un transcript RESUMABLE) AVANT de sortir. En control-group (défaut), systemd
# SIGKILLerait simultanément sudo→node→claude, tronquant le tour en cours → session
# non-relançable. mixed laisse Atelier piloter l'arrêt propre, puis SIGKILL le reste du
# cgroup seulement après TimeoutStopSec (> budget de drain ATELIER_AGENT_DRAIN_SECS=45s).
KillMode=mixed
TimeoutStopSec=60s

# Security hardening
ProtectSystem=strict
# /var/lib/hr-studio : home de l'utilisateur hr-studio, où l'agent runner (lancé en
# hr-studio via sudo) fait persister les sessions du Claude Agent SDK
# (CLAUDE_CONFIG_DIR=/var/lib/hr-studio/.claude). ProtectSystem=strict rend tout le FS
# read-only hors ReadWritePaths, et le namespace mount est hérité par les descendants
# (sudo→node→claude) : sans cette entrée, la CLI échoue SILENCIEUSEMENT à écrire le
# transcript (EROFS) → pas de reprise/historique de conversation. Cf. multi-conversations.
ReadWritePaths=/run /var/lib/atelier /var/lib/hr-studio /opt/atelier/data /tmp
ProtectHome=false

[Install]
WantedBy=multi-user.target
