# systemd/atelier-sync-docs.service
[Unit]
Description=Atelier — sync docs from Medion (source of truth)
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/atelier/bin/sync-docs.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atelier-sync-docs
