# systemd/atelier-sync-git.service
[Unit]
Description=Atelier — sync git bare repos from Medion (source of truth)
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/atelier/bin/sync-git.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atelier-sync-git
