# systemd/atelier-sync-store.service
[Unit]
Description=Atelier — sync app store from Medion (source of truth)
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/atelier/bin/sync-store.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atelier-sync-store
