# systemd/atelier.service
[Unit]
Description=Atelier — Application Platform
# Le control-plane (registres apps+ports, tâches, index docs) vit désormais dans
# Postgres (atelier_meta). La supervision en dépend au démarrage — on ordonne donc
# après postgresql. Pas de Requires= : si Postgres tombe à chaud, Atelier continue
# de tourner (les écritures dégradent), mais un boot sans Postgres échoue volontairement.
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/atelier/bin/atelier
WorkingDirectory=/opt/atelier
EnvironmentFile=-/opt/atelier/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atelier

# Security hardening
ProtectSystem=strict
ReadWritePaths=/run /var/lib/atelier /opt/atelier/data /tmp
ProtectHome=false

[Install]
WantedBy=multi-user.target
