# systemd/atelier-sync-docs.timer
[Unit]
Description=Atelier — sync docs every 5 minutes
Requires=atelier-sync-docs.service

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
AccuracySec=10s
Unit=atelier-sync-docs.service

[Install]
WantedBy=timers.target
