# systemd/atelier-sync-store.timer
[Unit]
Description=Atelier — sync app store every 5 minutes
Requires=atelier-sync-store.service

[Timer]
OnBootSec=45s
OnUnitActiveSec=5min
AccuracySec=10s
Unit=atelier-sync-store.service

[Install]
WantedBy=timers.target
