# systemd/hr-flowd.service
[Unit]
Description=hr-flowd — HomeRoute flow daemon (multi-stack, callback-based)
After=network-online.target
Wants=network-online.target
# Loosely coupled to atelier.service: hr-flowd doesn't require it to start,
# but a healthy Atelier API is what fronts the daemon for write traffic.
PartOf=atelier.service

[Service]
Type=simple
ExecStart=/opt/atelier/bin/atelier-flowd
WorkingDirectory=/opt/atelier
EnvironmentFile=-/etc/atelier-flowd/.env
EnvironmentFile=-/opt/atelier/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hr-flowd

# Reload TOML registry on SIGHUP (cf. Phase 1.9). Equivalent to:
#   curl -H "Authorization: Bearer $TOK" -X POST localhost:4002/v1/_admin/reload
ExecReload=/bin/kill -HUP $MAINPID

# Security hardening — same shape as atelier.service.
ProtectSystem=strict
ReadWritePaths=/run /var/lib/atelier /opt/atelier/data /tmp
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
